# buildroot_external/package/gvforwarder/systemd/gvforwader.service
[Unit]
Description=gvisor-tap-vsock Network Traffic Forwarder
Requires=podman.service
After=podman.service
StartLimitIntervalSec=0
ConditionVirtualization=!wsl

[Service]
Type=simple
Restart=on-failure
ExecStop=/bin/kill -s TERM $MAINPID
ExecStart=/usr/bin/gvforwarder

[Install]
WantedBy=default.target
